# snapwebsites/snapinit/debian/snapinit.service
[Unit]
Description=Snap! Websites snapinit daemon
After=network.target

[Service]
Type=forking
EnvironmentFile=/etc/default/snapinit
ExecStart=/etc/init.d/snapinit start
ExecStop=/etc/init.d/snapinit stop
Restart=on-failure
PIDFile=/run/lock/snapwebsites/snapinit-lock.pid

#[Install]
#WantedBy=multi-user.target
