# snapwebsites/snapbase/debian/snapbase.cron.monthly
#!/bin/sh -e
#

# Delete old files under /var/cache/snapwebsites
#
find /var/cache/snapwebsites -type f -atime +31 -delete

