# snapwebsites/snapbase/debian/snapbase.cron.monthly
#!/bin/sh -e
#

# Delete old files and empty directories under /var/cache/snapwebsites
#
find /var/cache/snapwebsites/* -type f -atime +31 -delete
find /var/cache/snapwebsites/* -type d -exec rmdir --ignore-fail-on-non-empty "{}" \;

